# mac_transcriber/launchd/com.slack-zoom.gigaam-transcriber.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.slack-zoom.gigaam-transcriber</string>

  <key>WorkingDirectory</key>
  <string>/Users/ilyaqa/Projects/GigaAM</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/ilyaqa/Projects/GigaAM/.venv/bin/python</string>
    <string>-m</string>
    <string>uvicorn</string>
    <string>mac_transcriber.service:app</string>
    <string>--host</string>
    <string>127.0.0.1</string>
    <string>--port</string>
    <string>18003</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>MPLCONFIGDIR</key>
    <string>/tmp/slack-zoom-matplotlib</string>
    <key>HF_TOKEN</key>
    <string></string>
    <key>MAC_TRANSCRIBER_API_KEY</key>
    <string>REPLACE_WITH_SHARED_SECRET</string>
    <key>MAC_TRANSCRIBER_ROOT</key>
    <string>/Users/ilyaqa/Projects/GigaAM/.local/mac_transcriber</string>
    <key>MAC_TRANSCRIBER_MODEL</key>
    <string>v3_e2e_rnnt</string>
    <key>MAC_TRANSCRIBER_CACHE</key>
    <string>/tmp/gigaam-cache</string>
    <key>MAC_TRANSCRIBER_DEVICE</key>
    <string>cpu</string>
    <key>MAC_TRANSCRIBER_DIARIZATION_DEVICE</key>
    <string>cpu</string>
    <key>MAC_TRANSCRIBER_BATCH_SIZE</key>
    <string>4</string>
    <key>MAC_TRANSCRIBER_REPORT_MODE</key>
    <string>ai</string>
    <key>MAC_TRANSCRIBER_REPORT_MODEL</key>
    <string>gpt-5.5</string>
    <key>MAC_TRANSCRIBER_BASELINE_UPGRADE_MODEL</key>
    <string></string>
    <key>MAC_TRANSCRIBER_REPORT_PDF</key>
    <string>1</string>
    <key>MAC_TRANSCRIBER_AI_SYNTHESIS_CHUNK_LIMIT</key>
    <string>24</string>
    <key>MAC_TRANSCRIBER_AI_SYNTHESIS_BATCH_SIZE</key>
    <string>4</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/ilyaqa/Library/Logs/slack_zoom/gigaam-transcriber.out.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/ilyaqa/Library/Logs/slack_zoom/gigaam-transcriber.err.log</string>
</dict>
</plist>
